# linux_files/rc.local
#!/bin/sh -e
# Start Home Automation
/usr/bin/python /home/pi/rpi-home-automation/master_controller.py > /dev/null 2>&1
exit 0
